# etc/crontabs/abc
0 6 * * * /etc/s6-overlay/scripts/cleanup-database
