# etc/crontabs/abc
0 6 * * * /usr/bin/sqlite3 /app/.ehforwarderbot/profiles/default/blueset.telegram/tgdata.db "DELETE FROM msglog WHERE time < date('now','-7 day'); VACUUM"
