# etc/services.d/mcl/run
#!/usr/bin/with-contenv bash

echo "Starting mcl."

cd /app/mcl

exec s6-setuidgid abc ./mcl -g "--no-console"
