# etc/s6-overlay/s6-rc.d/cron/up
#!/usr/bin/with-contenv bash

crontab -u abc /etc/crontabs/abc
