# etc/s6-overlay/s6-rc.d/bot/run
#!/usr/bin/with-contenv bash

echo "Starting bot."

exec s6-setuidgid abc python -m ehforwarderbot --profile $EFB_PROFILE
