# etc/s6-overlay/s6-rc.d/bot/up
#!/usr/bin/with-contenv bash

if [ ! -d "$EFB_DATA_PATH" ]; then
    mkdir -v "$EFB_DATA_PATH"
fi
chown -R abc:abc "$EFB_DATA_PATH"
